# revd
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
PROG=/usr/bin/python3
SCRIPT_PATH=$ROOT_DIR/bot_openwrt.py
LOG_FILE=/var/log/revd_bot.log

start_service() {
    # Check if script exists
    if [ ! -f "\$SCRIPT_PATH" ]; then
        echo "Skrip bot tidak ditemukan di \$SCRIPT_PATH"
        logger -t revd "Skrip bot tidak ditemukan di \$SCRIPT_PATH"
        return 1
    fi
    
    # Check if Python3 is installed
    if [ ! -x "\$PROG" ]; then
        echo "Python3 tidak ditemukan di \$PROG"
        logger -t revd "Python3 tidak ditemukan di \$PROG"
        return 1
    fi

    # Make sure plugins directory exists and has correct permissions
    if [ ! -d "$PLUGINS_DIR" ]; then
        mkdir -p $PLUGINS_DIR
        chmod 755 $PLUGINS_DIR
        logger -t revd "Membuat direktori plugins"
    fi

    # Check for plugin scripts and copy from backup if missing
    for script in speedtest.sh reboot.sh ping.sh clear_ram.sh vnstat.sh system.sh; do
        if [ ! -f "$PLUGINS_DIR/\$script" ] && [ -f "/etc/revd_backup/\$script" ]; then
            cp "/etc/revd_backup/\$script" "$PLUGINS_DIR/\$script"
            chmod +x "$PLUGINS_DIR/\$script"
            logger -t revd "Memulihkan skrip \$script dari backup"
        fi
    done
    
    # Log starting message
    logger -t revd "Memulai layanan Telegram Bot"
    
    # Configure the service with logging
    procd_open_instance
    procd_set_param command \$PROG \$SCRIPT_PATH
    procd_set_param stderr 1
    procd_set_param stdout 1
    procd_set_param respawn \${respawn_threshold:-3600} \${respawn_timeout:-5} \${respawn_retry:-5}
    procd_close_instance
    
    # Create a log entry for successful start
    echo "\$(date): Service started" >> \$LOG_FILE
}

stop_service() {
    # Log stopping message
    logger -t revd "Menghentikan layanan Telegram Bot"
    
    # Find and kill all Python processes running the bot script
    PIDS=\$(ps | grep "\$SCRIPT_PATH" | grep -v grep | awk '{print \$1}')
    if [ -n "\$PIDS" ]; then
        kill -9 \$PIDS 2>/dev/null
        logger -t revd "Bot dihentikan, PID: \$PIDS"
    else
        logger -t revd "Tidak ada proses bot yang berjalan"
    fi
    
    # Create a log entry for stop
    echo "\$(date): Service stopped" >> \$LOG_FILE
}

reload_service() {
    logger -t revd "Memuat ulang layanan Telegram Bot"
    stop
    sleep 2
    start
}

restart() {
    logger -t revd "Memulai ulang layanan Telegram Bot"
    stop
    sleep 3
    start
}
